# ifaceDown/citiIfaceDown.conf
start on started nwdecoder
exec /usr/sbin/citTest.sh